# scripts/com.sarahbroussard.portfolio-autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sarahbroussard.portfolio-autopush</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sarahbroussard/CascadeProjects/sarah-portfolioV2/scripts/auto-push.sh</string>
    </array>

    <!-- Run every hour (3600 seconds) -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Also run once 1 minute after load -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/sarahbroussard/CascadeProjects/sarah-portfolioV2/.auto-push-logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sarahbroussard/CascadeProjects/sarah-portfolioV2/.auto-push-logs/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/sarahbroussard</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/sarahbroussard/CascadeProjects/sarah-portfolioV2</string>
</dict>
</plist>
